# package/independent/base-settings-xwrt-xf-949/files/rc.local

env -i ACTION="ifchange" INTERFACE="br-lan" LINK=1 /sbin/hotplug-call iface

exit 0
